# root/etc/s6-overlay/s6-rc.d/init-rsnapshot-upgrade/run
#!/usr/bin/with-contenv bash

echo "***** Upgrading config file if required *****"

rsnapshot upgrade-config-file
