# root/etc/s6-overlay/s6-rc.d/init-rsnapshot-test/run
#!/usr/bin/with-contenv bash

echo "***** Validating config file *****"

rsnapshot configtest
