# root/etc/s6-overlay/s6-rc.d/init-rsnapshot-test/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "***** Validating config file *****"

s6-setuidgid abc rsnapshot configtest
